# services/linux/claude-observability.service
[Unit]
Description=Claude Code Observability Server
Documentation=https://github.com/gsnaiper/claude-code-hooks-multi-agent-observability
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.claude/observability/server
ExecStart=/usr/local/bin/bun run src/index.ts
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="NODE_ENV=production"
Environment="PORT=4000"

[Install]
WantedBy=default.target
